# deb/etc/cron.daily/wcfg-backup
#!/bin/sh

test -x /usr/bin/wcfg-backup || exit 0
/usr/bin/wcfg-backup

